# speed_test.timer
# path:   /home/klassiker/.local/share/repos/shell/speed_test.timer
# author: klassiker [mrdotx]
# github: https://github.com/mrdotx/shell
# date:   2024-04-19T10:12:01+0200

[Unit]
Description=run speed_test.service twice a day

[Timer]
OnCalendar=9,21:37:00
AccuracySec=1m

[Install]
WantedBy=timers.target
