# systemd/pkgstats.timer
# path:   /home/klassiker/.local/share/repos/shell/systemd/pkgstats.timer
# author: klassiker [mrdotx]
# github: https://github.com/mrdotx/shell
# date:   2024-04-30T08:27:17+0200

[Unit]
Description=run pkgstats.service monthly

[Timer]
OnCalendar=*-*-02 04:04:00
AccuracySec=1m

[Install]
WantedBy=timers.target
